# trunk/falcony/laptop-mode-tools/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/acpi/
drwxr-xr-x	root/root	etc/acpi/actions/
-rwxr-xr-x	root/root	etc/acpi/actions/lm_ac_adapter.sh
-rwxr-xr-x	root/root	etc/acpi/actions/lm_battery.sh
-rwxr-xr-x	root/root	etc/acpi/actions/lm_lid.sh
drwxr-xr-x	root/root	etc/acpi/events/
-rwxr-xr-x	root/root	etc/acpi/events/lm_ac_adapter
-rwxr-xr-x	root/root	etc/acpi/events/lm_battery
-rwxr-xr-x	root/root	etc/acpi/events/lm_lid
drwxr-xr-x	root/root	etc/laptop-mode/
-rwxr-xr-x	root/root	etc/laptop-mode/laptop-mode.conf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/laptop-mode
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man8/
-rw-r--r--	root/root	usr/man/man8/laptop-mode.conf.8.gz
-rw-r--r--	root/root	usr/man/man8/laptop_mode.8.gz
drwxr-xr-x	root/root	usr/sbin/
-rwx------	root/root	usr/sbin/laptop_mode
